# dot_config/systemd/user/xremap.service
[Unit]
Description=xremap

[Service]
Restart=always
ExecStart=%h/bin/xremap %h/.config/xremap.yml --watch
